# crontab/finishing.crontab
@reboot sudo passwd -l root
@reboot sudo passwd -l pi
@reboot sudo crontab -r
